# nas_ports/freenas/freenas-serviced/pkg-plist
/usr/local/etc/serviced.d/org.freebsd.rc.json
/usr/local/etc/rc.d/serviced
/usr/local/sbin/serviced
/usr/local/sbin/servicectl
