# nas_ports/freenas/freenas-crashd/pkg-plist
etc/rc.d/crashd
sbin/crashd
sbin/crash-wrapper
