# nas_ports/freenas/freenas-containerd/pkg-plist
etc/rc.d/containerd
sbin/containerd
share/containerd/virtio-win.iso
